# verium.conf
[program:vericoin]
process_name=verium
command=/home/ubuntu/verium-master/cpuminer-opt-3.6.7/run.sh
autostart=true
autorestart=true
stderr_logfile=/home/ubuntu/verium-master/log.log
